# deploy/launchd/com.adrimg3196.lucidfence.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.adrimg3196.lucidfence</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>python3</string>
    <string>saas_server.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Applications/LucidFence.app/Contents/Resources</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>LUCIDFENCE_PORT</key>
    <string>8765</string>
    <key>LUCIDFENCE_HOST</key>
    <string>127.0.0.1</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/lucidfence-launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/lucidfence-launchd.err.log</string>
</dict>
</plist>
